# control-loop/adaptation/supervisor-app.conf
[supervisord]
loglevel = debug

[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /app/uwsgi.ini
user=app

[program:nginx-app]
command = /usr/sbin/nginx
loglevel = debug
